# runit/sv/clapboard/run
#!/bin/sh
exec 2>&1
exec clapboard --record "both" 2>&1
