# runit/sv/wpaperd/run
#!/bin/sh
[ -r ./conf ] && ./conf
exec /usr/bin/env -S wpaperd 2>&1
# vim:ft=bash
