# runit/sv/ssh-agent/run
#!/bin/sh
exec 2>&1
exec ssh-agent -D -a "${XDG_RUNTIME_DIR}/ssh-agent.socket" 2>&1
