# runit/sv/dunst/run
#!/bin/sh
exec /usr/bin/env -S dunst 2>&1
# vim:ft=bash
