# runit/sv/mpd/run
#!/bin/sh
[ -r ./conf ] && . ./conf
exec /bin/env -S mpd ${OPTS:---no-daemon} 2>&1
