# runit/sv/pipewire/run
#!/bin/sh
exec /usr/bin/env XDG_RUNTIME_DIR="/run/user/$(id -u)" pipewire 2>&1
