# runit/sv/polkit-agent/run
#!/bin/sh
exec /usr/bin/env -S polkit-agent 2>&1
# vim:ft=bash
